# profile.d/infinality.sh
/etc/fonts/infinality/rendering-settings.sh